# deploy/systemd/imgd.service
[Unit]
Description=imgd upload service (axum)
After=network.target

[Service]
Type=simple
User=imgd
Group=imgd
WorkingDirectory=/opt/imgd
ExecStart=/opt/imgd/imgd
Environment=PORT=3000
Environment=UPLOAD_TOKEN=replace-with-long-random-token
Environment=PUBLIC_BASE_URL=https://img.example.com/images
Environment=DATA_DIR=/data/images
Environment=RUST_LOG=imgd=info,tower_http=info
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/images
UMask=0027
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
